# recycler.service
[Unit]
Description=Clean Recycler

[Service]
Type=oneshot
ExecStart=/usr/bin/recycler -e

[Install]
WantedBy=multi-user.target
