# pkg/recycler/usr/lib/systemd/system/recycler.service
[Unit]
Description=Clean Recycler

[Service]
ExecStart=/usr/bin/recycler -e

[Install]
WantedBy=multi-user.target
